# basketballshop/geekshop.service
[Unit]
Description=geekshop daemon
After=network.target

[Service]
User=jo
Group=www-data
WorkingDirectory=/home/jo/geekshop
ExecStart=/usr/bin/gunicorn3 --access-logfile server.logs --workers 2 --bind unix:/home/jo/geekshop/geekshop.sock geekshop.wsgi

[Install]
WantedBy=multi-user.target